# services/mysql.service.conf
[program:mysql]
user=root
command=/usr/bin/pidproxy /var/run/mariadb/mysqld.pid /usr/bin/mysqld_safe --pid-file=/var/run/mariadb/mysqld.pid
stdout_events_enabled = true
stderr_events_enabled = true
